# systemd/fleetsign.service
# Started by ~/.config/labwc/autostart (see install.sh), not by `systemctl enable`:
# the desktop's graphical-session.target is not reliably reached for --user units on
# Raspberry Pi OS. systemd just supervises it (Restart=always).
[Unit]
Description=FleetSign

[Service]
Type=simple
WorkingDirectory=%h/fleetsign
ExecStart=%h/fleetsign/.venv/bin/python -m fleetsign --root %h/fleetsign
Restart=always
RestartSec=3
# Tag journal entries so `journalctl --user -t fleetsign` works alongside -u.
SyslogIdentifier=fleetsign
# Log verbosity is INFO by default. For verbose troubleshooting, uncomment the
# next line (or run `systemctl --user edit fleetsign` to add it as an override),
# then `systemctl --user restart fleetsign`. Levels: debug, info, warning, error.
#Environment=FLEETSIGN_LOG_LEVEL=debug
